# deploy/oci/systemd/vex-ranker-backup.timer
[Unit]
Description=Nightly VEX Ranker Postgres backup

[Timer]
OnCalendar=*-*-* 03:20:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
