# legacy-python-monitor/ops/com.vexranker.monitor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.vexranker.monitor</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/youruser/Documents/Projects/Vex Ranker/.venv/bin/python</string>
      <string>/Users/youruser/Documents/Projects/Vex Ranker/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/youruser/Documents/Projects/Vex Ranker</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/youruser/Documents/Projects/Vex Ranker/logs/monitor-service.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/youruser/Documents/Projects/Vex Ranker/logs/monitor-service.err.log</string>
  </dict>
</plist>
